# contrib/installer/debian/postinst.in
#!/bin/sh

TAO_VERSION=2.1.8

TAO_FILES="libTAO.so \
    libTAO_Utils.so \
    libTAO_PI.so \
    libTAO_AnyTypeCode.so \
    libTAO_PortableServer.so \
    libTAO_CodecFactory.so \
    libTAO_Codeset.so"

if [ "$1" = "configure" ]; then
    for i in $TAO_FILES
    do
      ln -s /opt/lib/$i.$TAO_VERSION /opt/lib/$i.1
    done
    ln -s /opt/lib/libTAO_PI.so.$TAO_VERSION /opt/lib/libTAO_PI.so
    ln -s /opt/lib/libTAO_Codeset.so.$TAO_VERSION /opt/lib/libTAO_Codeset.so
    ln -s /opt/lib/libtofdll.so.1.0.0 /opt/lib/libtofdll.so.1
    ln -s /opt/lib/libtofdll.so.1.0.0 /opt/lib/libtofdll.so

    echo "/opt/lib" > /etc/ld.so.conf.d/hmqd.conf
    ldconfig

    chmod 755 /etc/init.d/hmqd
    update-rc.d hmqd defaults
fi
